# crontab
*/30 * * * * cd /path/to/vendor-nav; source .source.bash; date > log; ./scrape.bash > log
